# etc/systemd/system/systemd-timesyncd.service.d/candle.conf
[Unit]
ConditionPathExists=!/boot/candle_hardware_clock.txt
