# etc/systemd/system/bthelper@.service.d/fix.conf
[Unit]
Before=

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 2
