# etc/systemd/system/seeed-voicecard.service.d/candle.conf
[Unit]
ConditionPathExists=!/boot/firmware/do_not_use_repeaker_hat.txt
ConditionPathExists=/home/pi/candle/respeaker_check.sh

[Service]
ExecCondition=/bin/bash /home/pi/candle/respeaker_check.sh
