# etc/systemd/system/NetworkManager.service.d/candle.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/NetworkManager --no-daemon --intern-config /home/pi/.webthings/etc/NetworkManager/intern-config
