# etc/cron.d/every_minute
* * * * * root /home/pi/candle/every_minute.sh
