# etc/cron.d/every_hour
0 * * * * root /home/pi/candle/every_hour.sh
