# etc/cron.d/every_minute_user
* * * * * pi /home/pi/candle/every_minute_user.sh
