# etc/systemd/system/rc-local.service.d/candle.conf
[Unit]
After=network.target candle_early.service

[Service]
ExecStartPre=/bin/chmod +x /etc/rc.local
