# etc/systemd/system/candle_splash_video.service
[Unit]
Description=Show Candle video logo during boot
After=graphical.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/candle
ExecStart=/home/pi/candle/flame.sh            
StandardOutput=journal+console
StandardError=journal+console


[Install]
WantedBy=multi-user.target
