# etc/systemd/system/candle_every_minute.timer
[Unit]
Description=Candle every minute timer
After=candle_late.service

[Timer]
OnCalendar=*-*-* *:*:00
AccuracySec=5s
Persistent=true
Unit=candle_every_minute.service

[Install]
WantedBy=timers.target
