# etc/systemd/system/triggerhappy.service.d/candle.conf
[Service]
User=root
ExecStart=/usr/sbin/thd --triggers /etc/triggerhappy/triggers.d/ --socket /run/thd.socket --user root --deviceglob /dev/input/event*
