# etc/systemd/system/dbus-fi.w1.wpa_supplicant1.service.d/candle.conf
[Service]
#ExecStart=/usr/sbin/wpa_supplicant -u -O "DIR=/run/wpa_supplicant GROUP=netdev" 
#ExecStart/usr/sbin/wpa_supplicant -d -i wlan0 -c /etc/wpa_supplicant/wpa_supplicant.conf -u p2p_set disabled 1 -N -i uap0 -c /etc/wpa_supplicant/wpa_supplicant.conf -u p2p_set disabled 1 -N -i wlan1 -u -O "DIR=/run/wpa_supplicant GROUP=netdev" 
ExecStart=
ExecStart=/bin/bash -c 'if ip link show | grep -q "wlan1:"; then exec /usr/sbin/wpa_supplicant -d -i wlan0 -c /etc/wpa_supplicant/wpa_supplicant.conf -u p2p_set disabled 1 -N -i uap0 -c /etc/wpa_supplicant/wpa_supplicant.conf -u p2p_set disabled 1 -N -i wlan1 -u -C "DIR=/run/wpa_supplicant GROUP=netdev" ; else exec /usr/sbin/wpa_supplicant -d -i wlan0 -c /etc/wpa_supplicant/wpa_supplicant.conf -u p2p_set disabled 1 -N -i uap0 -c /etc/wpa_supplicant/wpa_supplicant.conf -u p2p_set disabled 1; fi'
