# etc/systemd/system/dbus-fi.w1.wpa_supplicant1.service.d/candle.conf
[Service]
#ExecStart=/usr/sbin/wpa_supplicant -u -O "DIR=/run/wpa_supplicant GROUP=netdev" 
ExecStart=
ExecStart/usr/sbin/wpa_supplicant -c /etc/wpa_supplicant/wpa_supplicant.conf -u
