# etc/systemd/system/rsyslog.service.d/candle.conf
[Unit]
ConditionPathExists=/boot/firmware/developer.txt
