# etc/systemd/system/wpa_supplicant.service.d/candle.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/wpa_supplicant -c /etc/wpa_supplicant/wpa_supplicant.conf -iwlan0 -u -s -O "DIR=/run/wpa_supplicant GROUP=netdev"
