# etc/systemd/system/bluetooth.service.d/candle.conf
[Unit]
After=home-pi-.webthings.mount

[Service]
ExecStart=
ExecStart=/usr/libexec/bluetooth/bluetoothd --noplugin=sap
